# config/includes.chroot/etc/systemd/system/ssh.service.d/nomad.conf
[Unit]
# sshd must start AFTER nomad-ssh-keys has copied the persisted host
# keys into place. Without this ordering, sshd could come up with
# the fresh per-boot keys and users would get "REMOTE HOST KEY
# CHANGED" warnings.
#
# sshd must ALSO wait for nomad-passwords which sets the nomad user's
# password to "nomad". Without this, the first SSH connection of the
# boot finds the live-config default password ("live") and rejects
# what the MOTD told the user to try.
#
# We use After=, not Requires=, deliberately: if these setup services
# fail (e.g. data partition missing, which would already have broken
# everything), we still want sshd to come up with whatever state it
# has, so the user can at least log in and diagnose.
After=nomad-ssh-keys.service nomad-passwords.service
