# config/includes.chroot/etc/systemd/system/nomad-mount.service
[Unit]
Description=Mount NOMAD_DATA partition and set up data-partition state
# Everything nomad-related depends on this. Docker's drop-in adds
# Requires/After us, so we don't list docker here.
Before=docker.service nomad-init.service nomad-start.service nomad-ssh-keys.service

[Service]
Type=oneshot
RemainAfterExit=yes

# The actual work is in a script because there's enough of it that
# inline ExecStart= commands would be unreadable and hard to reason
# about error handling on.
ExecStart=/usr/local/sbin/nomad-mount

# Clean unmount on shutdown — tears down symlink, flushes writes.
ExecStop=/usr/local/sbin/nomad-mount stop

[Install]
WantedBy=multi-user.target
