# config/includes.chroot/etc/systemd/system/nomad-passwords.service
[Unit]
Description=Nomad password setup (every boot — live overlay resets /etc/shadow)
# Run BEFORE ssh so the first SSH connection attempt finds the
# correct password. Previously this ran After=multi-user.target
# which meant sshd could come up with the default live-config
# password ("live") and reject nomad/nomad logins for a window.
#
# We still need user-setup to have created the nomad user. That's
# part of live-config which hooks in at sysinit.target. We depend
# on sysinit.target being complete, but run before ssh starts.
After=sysinit.target
Before=ssh.service

# Skip if we've already run this boot. We use /run (tmpfs) for the
# flag so it's fresh every boot — this service should run exactly
# once per boot, not once per USB, because the password gets wiped
# along with /etc/shadow when the RAM overlay resets.
ConditionPathExists=!/run/nomad-passwords.done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/nomad-passwords
ExecStartPost=/bin/touch /run/nomad-passwords.done

[Install]
WantedBy=multi-user.target nomad-firstboot.target
