# config/includes.chroot/etc/systemd/system/docker.service.d/nomad.conf
[Unit]
# Docker must not start before NOMAD_DATA is mounted. Otherwise docker
# would create /mnt/nomad-data/docker on the unmounted mount point,
# start there, and then get shadowed when the real mount happens.
Requires=nomad-mount.service
After=nomad-mount.service
