# config/includes.chroot/etc/systemd/system/getty@tty1.service.d/nomad-clear.conf
[Unit]
# Wait until our nomad services have stabilized before we put a
# login prompt on tty1. Without this, getty races against
# nomad-start: getty may grab tty1 first, draw the issue/login
# prompt, then nomad-start (or any other service printing to console)
# scribbles over it. Looks confusing.
After=nomad-start.service

[Service]
# These are inherited from the upstream getty@.service template, but
# we set them explicitly here as a defensive measure: TTYReset and
# TTYVHangup ensure tty1 gets a fresh kernel-level reset before
# /sbin/agetty paints the login prompt, wiping any stray output left
# behind by earlier services (like live-boot's rsync).
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Belt-and-suspenders: ExecStartPre runs `clear` against the tty
# *before* agetty starts. tput clear works without TERM set, which
# we don't have at this point in boot. The `< /dev/tty1` redirect
# attaches input from tty1 so tput knows what terminal to clear.
ExecStartPre=-/usr/bin/clear
